# supervisord.conf
[supervisord]
logfile=./app-server.log   
loglevel=info

[program:prospectiva]
command=python3.8 main.py